# server/com.woojaedabin.wedding-backup.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.woojaedabin.wedding-backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/dbkim/woojae-dabin-wedding-server/backup_sqlite.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/dbkim/woojae-dabin-wedding-server</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/Users/dbkim/woojae-dabin-wedding-server/logs/backup.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dbkim/woojae-dabin-wedding-server/logs/backup.err.log</string>
</dict>
</plist>
